# gwdg-cloud-share-server/DEBIAN/prerm
#! /bin/bash
# prerm script

set -e


if [ -x "/etc/init.d/powerfolder-server" ]; then
        invoke-rc.d powerfolder-server stop || exit $?
fi

exit 0
